# system/getty@tty1.service.d/activate-numlock.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty '-p -- \\u' --nohints --noclear %I $TERM
